# profiles/homeserver/modules/homeserver-update-runtime/files_root/etc/systemd/system/harmonia-homeserver.service
[Unit]
Description=Run Harmonia HOMESERVER rolling-release convergence
Documentation=file:/var/lib/harmonia/receipts/homeserver-update-latest/run.json
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/harmonia homeserver-update /etc/harmonia/profiles/homeserver/index.json --apply --receipt-dir /var/lib/harmonia/receipts/homeserver-update-latest
Nice=10
IOSchedulingClass=idle
