# profiles/homeconsole/modules/homeconsole-update-runtime/files_root/etc/systemd/system/harmonia-homeconsole.service
[Unit]
Description=Run Harmonia HomeConsole profile convergence
Documentation=file:/var/lib/harmonia/receipts/homeconsole-update-latest/run.json
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/harmonia homeconsole-update /etc/harmonia/profiles/homeconsole/index.json --apply --receipt-dir /var/lib/harmonia/receipts/homeconsole-update-latest
Nice=10
IOSchedulingClass=idle
