# profiles/homeconsole/modules/homeconsole-update-runtime/files_root/etc/systemd/system/harmonia-homeconsole.timer
[Unit]
Description=Run Harmonia HomeConsole profile convergence on schedule

[Timer]
OnBootSec=2min
OnCalendar=*:0/10
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=harmonia-homeconsole.service

[Install]
WantedBy=timers.target
