# profiles/homeserver/modules/matrix/files_root/etc/systemd/system/harmonia-matrix-prepare.service
[Unit]
Description=Prepare Matrix Synapse product substrate
After=postgresql.service network.target
Wants=postgresql.service
Before=matrix-synapse.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/harmonia-matrix-converge
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
