# profiles/tv/modules/tv-update-runtime/files_root/etc/systemd/system/harmonia-tv.timer
[Unit]
Description=Run Harmonia TV rolling-release convergence on schedule

[Timer]
OnBootSec=2min
OnCalendar=*:0/10
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true
Unit=harmonia-tv.service

[Install]
WantedBy=timers.target