# profiles/tv/modules/tv-update-runtime/files_root/etc/systemd/system/harmonia-tv.service
[Unit]
Description=Run Harmonia TV rolling-release convergence
Documentation=file:/var/lib/harmonia/receipts/tv-update-latest/run.json
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/harmonia tv-update /etc/harmonia/profiles/tv/index.json --apply --receipt-dir /var/lib/harmonia/receipts/tv-update-latest
Nice=10
IOSchedulingClass=idle